# systemd/attendance.service
[Unit]
Description=AI Attendance Raspberry Pi Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/attendance-pi-client
ExecStart=/usr/bin/python3 /opt/attendance-pi-client/app.py
Restart=always
RestartSec=5
User=pi
Group=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

