# systemd/attendance.service
[Unit]
Description=AI Attendance Raspberry Pi Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/attendance/attendance-pi-client
ExecStart=/home/attendance/attendance-pi-client/venv/bin/python /home/attendance/attendance-pi-client/app.py
Restart=always
RestartSec=5
User=attendance
Group=attendance
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target